# wrt-go/files/etc/uci-defaults/50-vmtools

. /lib/functions/system.sh

service vmtoolsd disable
service qemu-ga disable

case "$(board_name)" in
    *vmware*)
        service vmtoolsd enable
        service vmtoolsd start
        ;;
    *)
        service qemu-ga enable
        service qemu-ga start
        ;;
esac

exit 0
